# nabd/nabd.socket
[Socket]
ListenStream=127.0.0.1:10543

[Install]
WantedBy = sockets.target
